# supervisord-playsmsd.conf
[program:playsmsd]
command=/start-playsmsd.sh
numprocs=1
autostart=true
autorestart=false
