# playsms/supervisord-php-fpm.conf
[program:php-fpm8.3]
command=/runner_php-fpm.sh
numprocs=1
autostart=true
autorestart=true
