# playsms/supervisord-playsmsd.conf
[program:playsmsd_runner]
command=/runner_playsmsd.sh
numprocs=1
autostart=true
autorestart=true
